# debian/sns-utils.postrm
#!/bin/sh -e

update-rc.d sns remove

#DEBHELPER#
