# debian/sns-utils.postinst
#!/bin/sh -e

update-rc.d sns defaults

#DEBHELPER#
